# systemd/marimo.service.template
[Unit]
Description=TPU Dev Marimo
After=network-online.target

[Service]
Type=simple
WorkingDirectory={{ENV_DIR}}
Environment=PATH={{VENV_DIR}}/bin:{{HOME}}/.local/bin:{{HOME}}/.cargo/bin:/usr/local/bin:/usr/bin:/bin
ExecStart={{VENV_DIR}}/bin/marimo edit --host 127.0.0.1 --port {{MARIMO_PORT}} --headless
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
